# deploy/systemd/seatseeker-db-backup.service
[Unit]
Description=SeatSeeker database backup job
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/SeatSeeker-UC-Merced-Course-Availability-Notifier
EnvironmentFile=-/etc/default/seatseeker-backup
ExecStart=/usr/bin/env bash /home/ubuntu/SeatSeeker-UC-Merced-Course-Availability-Notifier/deploy/scripts/backup_db.sh
